# roles/antivirus/templates/clamav-scan.timer.j2
# {{ ansible_managed }}
[Unit]
Description=ClamAV periodic virus scan
Documentation=man:clamscan(1)

[Timer]
OnCalendar={{ antivirus_scan_schedule }}
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
